# storyplayer/test-environments/vagrant-vbox-centos6-ssl/files/zmq/zmq-echo-server.conf
[program:zmq-echo-server]
command=/usr/local/bin/zmq-echo-server.php

